# rc.local
#!/bin/bash


cd /root ; ./audiophilize.sh &

exit 0
